# CMD-Path/PowerShellEnvMaintenance/EnvChecker.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <URI>\Ariel\PSEnvMaintenance</URI>
    <Description>PowerShell environment maintenance</Description>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2025-05-12T02:00:00-07:00</StartBoundary>
      <ScheduleByDay><DaysInterval>1</DaysInterval></ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Actions>
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-File "\\192.168.0.122\software\dev\Tools\Fix-PowerShellEnv.ps1"</Arguments>
    </Exec>
  </Actions>
</Task>
